# confs/rc.local
#!/bin/sh
# default rc.local script for ichirou
# See LICENSE file for license and copyright information.

echo
echo "----------------------------------"
echo "| ichirou init system - rc.local |"
echo "| place custom init scripts here |"
echo "----------------------------------"
echo
