# confs/rc.local.in
#!/bin/sh
# default rc.local script for ichirou
# Copyright (c) 2020-2021 Emily <elishikawa@jagudev.net>
# See LICENSE file or <https://www.gnu.org/licenses/gpl-3.0.html>
# for license and copyright information.

echo
echo "----------------------------------"
echo "| ichirou init system - rc.local |"
echo "| place custom init scripts here |"
echo "----------------------------------"
echo
